# mounts/media-WMG.mount
# Automatically generated by systemd-fstab-generator

[Unit]

[Mount]
Where=/media/WMG
What=//charite.de/centren/C17/WMG
Type=cifs
Options=credentials=/opt/.credentials_cifs,ro,nofail
